# deploy/systemd/simple-health-checker.service
Description=Simple Health Checker Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=penne
Group=penne
WorkingDirectory=/opt/simple_health_checker
EnvironmentFile=/opt/simple_health_checker/.env
Environment=PYTHONPATH=/opt/simple_health_checker/src
ExecStart=/opt/simple_health_checker/.venv/bin/python -m simple_health_checker.app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target